# init/roscore.conf
description "ROS Core starter"
author "info@teamdiana.org"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec roscore