# cron.d/daily-reboot
# Reboot the rover every day at 4

tell() {
wall << EOF
The rover will reboot in $1 hour(s) at 4:00, please save all work.
EOF
}

tellm() {
wall << EOF
The rover will reboot in $1 minute(s) at 4:00, please save all work.
EOF
}

* 1 * * * tell() 3 >/dev/null 2>&1
* 2 * * * tell() 2 >/dev/null 2>&1
* 3 * * * tell() 1 >/dev/null 2>&1
30 3 * * * tellm() 10 >/dev/null 2>&1
50 3 * * * tellm() 10 >/dev/null 2>&1
* 4 * * * reboot >/dev/null 2>&1